# packages/thorch-inputplumber/payload/usr/lib/systemd/system/inputplumber.service
[Unit]
Description=InputPlumber remaps input devices
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/inputplumber
Restart=always

[Install]
WantedBy=multi-user.target
