# packages/thorch-firstboot/payload/etc/skel/.config/autostart/thorch-firstboot.desktop
[Desktop Entry]
Type=Application
Name=Thorch Setup
Exec=/usr/bin/thorch-firstboot
OnlyShowIn=KDE;
Terminal=false
X-KDE-autostart-phase=2
